# service/alap.service
[Unit]
Description=Minden beallitas

[Service]
ExecStart=/home/bence/scriptek/alap.sh

[Install]
WantedBy=default.target
