# service/zene.service
[Unit]
Description=Zene stream
After=network.target

[Service]
ExecStart=/home/bence/scriptek/zene.sh

[Install]
WantedBy=default.target
